# templates/whelp.conf
[program:whelp]
command=/srv/whelp/whelp/straight_script/run_whelp.py
autostart=true
autorestart=true
user=www-data
